# systemd/kiosk-browser-user.service
[Unit]
Description=Office Kiosk Browser
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=%h/office_kiosk
ExecStart=%h/office_kiosk/start_kiosk.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
